# inst/sdb.crontab
# crontab copy

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# ---- [external db updates] ----

# ARGOS (every 3 hours)
0 */3       * * *  Rscript --default-packages=dup -e 'ARGOS.pipeline();'

# ARGOS2 (every hour)
0 */1       * * *  Rscript --default-packages=dup -e 'ARGOS2.pipeline();'


# RUFFatSEEWIESEN
# 0 18 * * *  Rscript --default-packages=dup -e 'RUFFatSEEWIESEN_photos.pipeline()'


#DRUID (every 3 hours)
0 */3 * * * Rscript --default-packages=dup -e  'DRUID.downloadNew(what = "GPS", logfile =  "log_druid_GPS.csv");DRUID.downloadNew(what = "ODBA", logfile =  "log_druid_ODBA.csv");DRUID.downloadNew(what = "ENV", , logfile =  "log_druid_ENV.csv")'



# ---- [Maintenance] ----

# replication
0 */3 * * *  Rscript --default-packages=dup -e 'check_replication("scidb", "scidb_maoraves")'

# backup
# 0 0        * * *    Rscript --default-packages=dup -e 'backup.pipeline()'     


# backup from field servers

# DUMMER

# 0 2 * * *  Rscript --default-packages=dup -e 'db_copy("FIELD_2025_NOLAatDUMMERSEE", "remotehost_1", "host")'     

#NZ
0 */1  * * *  Rscript --default-packages=dup -e 'db_copy("FIELD_2024_BADOatNZ", "remotehost_1", "host")'


# DOCUWIKI BACKUP
0 1        * * *    duplicity --no-encryption /media/DATA/WIKIS/html file:///ds/scidbbackup/dokuwiki_backups/duplicity/

# internal updates
3 22        * * *   Rscript --default-packages=dup -e 'DB_internal_updates.pipeline()'
