# inst/sdb_sudo.crontab
# crontab copy

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 2 * * * systemctl restart shiny-server
